# discord/second-brain-bot.service
[Unit]
Description=Second Brain Discord Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/second_brain
ExecStart=/usr/bin/python3 /home/pi/second_brain/discord/bot.py
Restart=always
RestartSec=10
EnvironmentFile=/home/pi/second_brain/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
